# home/Library/LaunchAgents/dev.cade.mlxserve.plist.tmpl
{{/*
LaunchAgent: keep mlx-openai-server running on localhost:8080.

This is what makes aider/opencode/pi "just work" when opened — without it,
each tool would fail to connect on first launch unless the user remembered
to run `mlxserve` first. Loaded by install/macos-services.sh.

Model: unsloth/Qwen3.6-27B-MLX-8bit (Apr 2026 release, 27B dense VLM,
reasoning-tuned, 256K native context). Served as "qwen3.6-27b" so configs
can reference a clean short name.

Quant: 8-bit weights. Recovers ~99.9% of FP vs ~98.9% for 4-bit —
important for diff-edit accuracy in agent loops.

Model-type note: Qwen3.6-27B is a VLM, but mlx-openai-server v1.8.1's VLM
batch scheduler is broken (always tries to pass kv_bits to BatchGenerator
which doesn't accept the kwarg — fails immediately on first request). Loading
the model under --model-type lm sidesteps the broken scheduler and uses the
LM batch path, which works for text-only requests (we don't need image input
for coding). Re-enable multimodal + kv-bits 8 once the VLM scheduler is fixed.

Cost: ~35GB RAM for weights + ~10-20GB KV at typical context lengths. On
128GB M3 Max that's plenty of headroom. To disable temporarily:
  launchctl bootout gui/$(id -u)/dev.cade.mlxserve   # or: mlxstop
*/ -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.cade.mlxserve</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/mlx-openai-server</string>
        <string>launch</string>
        <string>--model-type</string>
        <string>lm</string>
        <string>--model-path</string>
        <string>unsloth/Qwen3.6-27B-MLX-8bit</string>
        <string>--served-model-name</string>
        <string>qwen3.6-27b</string>
        <string>--tool-call-parser</string>
        <string>qwen3_coder</string>
        <string>--enable-auto-tool-choice</string>
        <string>--reasoning-parser</string>
        <string>qwen3_5</string>
        <string>--kv-bits</string>
        <string>8</string>
        <string>--kv-group-size</string>
        <string>64</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HF_HOME</key>
        <string>{{ .chezmoi.homeDir }}/.local/.cache/huggingface</string>
        <key>PATH</key>
        <string>{{ .chezmoi.homeDir }}/.local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/share/mlxserve/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/share/mlxserve/stderr.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
