# home/Library/LaunchAgents/dev.cade.cass-watch.plist.tmpl
{{/*
LaunchAgent: refresh cass's lexical index every five minutes. Persistent
`index --watch` can wedge in watch_startup while consuming multiple cores
(cass #258); bounded one-shot runs keep new sessions searchable without a
resident process. dev.cade.cass-semantic handles the expensive vectors daily.

cass is installed to ~/.local/bin by install/memory.sh (GitHub release binary,
checksum-verified). CASS_DATA_DIR keeps the SQLite index under ~/.cache
(scratch-linked on NFS machines).

Cost: no idle process; brief spikes while indexing new sessions. To disable:
  launchctl bootout gui/$(id -u)/dev.cade.cass-watch
*/ -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.cade.cass-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/cass</string>
        <string>index</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CASS_DATA_DIR</key>
        <string>{{ .chezmoi.homeDir }}/.cache/cass</string>
        <key>PATH</key>
        <string>{{ .chezmoi.homeDir }}/.local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/share/cass/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/share/cass/stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
